# supervisor.ini
[program:gin-template]                                                                                                                                          
command=/root/golang/gin-template/main -h 127.0.0.1 -p 7890 -debug=false -err=false -orm=false
directory=/root/golang/gin-template
autorestart=true
redirect_stderr=true
stdout_logfile=/root/log/gin-template.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
stopsignal=KILL
stopasgroup=true
